# services/couchpotato/run
#!/bin/bash

chmod 770 /etc/service/couchpotato/supervise
chown box:box -R /etc/service/couchpotato

exec setuser box python /opt/couchpotato/CouchPotato.py --data_dir /home/box/couchpotato