# deluge/services.d/deluge-web/run
#!/usr/bin/execlineb -P

s6-setuidgid box
cd /home/box/deluge/data
deluge-web -c config -L info -l config/deluge-web.log