# services/deluge-web/run
#!/bin/bash

chmod 770 /etc/service/deluge-web/supervise
chown box:box -R /etc/service/deluge-web/supervise

exec setuser box deluge-web -c /home/box/deluge -p 8112 --no-ssl