# services/deluge/run
#!/bin/bash

chmod 770 /etc/service/deluge/supervise
chown box:box -R /etc/service/deluge/supervise

exec setuser box deluged -c /home/box/deluge -d -L info -p 58846