# services/deluge/run
#!/bin/bash

chmod 770 /etc/service/deluged/supervise
chown box:box -R /etc/service/deluged/supervise

exec setuser box deluged -d -L info